# systemd/anpr_web.service
[Unit]
Description=ANPR Gate System — Web Dashboard
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/ANPR_System
ExecStart=/usr/bin/python3 -m web.app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
