# systemd/anpr_core.service
[Unit]
Description=ANPR Gate System — Core Recognition Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/ANPR_System
ExecStart=/usr/bin/python3 -m src.main --config /home/pi/ANPR_System/config.yaml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
